# acme-sh_beta/rootfs/etc/services.d/crond/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash

# Run crond
bashio::log.info "Starting crond..."
exec crond -f -L /dev/stdout
